# examples/launchagents/com.cyclisme.migration.10-phase1-install-now.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <!-- ═══════════════════════════════════════════════════════════════════ -->
    <!-- MIGRATION PHASE 1: Installation Nouveaux Agents -->
    <!-- ═══════════════════════════════════════════════════════════════════ -->

    <key>Label</key>
    <string>com.cyclisme.migration.10-phase1-install-now</string>

    <key>Description</key>
    <string>Migration Phase 1: Install new LaunchAgents with new naming convention</string>

    <!-- Execute immediately on load -->
    <key>RunAtLoad</key>
    <true/>

    <!-- Command -->
    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/stephanejouve/cyclisme-training-logs/scripts/launchagents/phase1-install-new-agents.sh</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/stephanejouve/cyclisme-training-logs/scripts/launchagents</string>

    <!-- Logging -->
    <key>StandardOutPath</key>
    <string>/Users/stephanejouve/Library/Logs/cyclisme-migration-phase1.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/stephanejouve/Library/Logs/cyclisme-migration-phase1.error.log</string>

    <!-- Environment -->
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
    </dict>
</dict>
</plist>
